# files/systemd/backup.timer
[Unit]
Description=Hourly home directory backups

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
