# files/systemd/backup.service
[Unit]
Description=Backup home directory

[Service]
Type=simple
User=alarm
ExecStart=/home/alarm/bin/backup.sh

[Install]
WantedBy=multi-user.target
