# files/systemd/backup.timer
[Unit]
Description=Quarter-hourly home directory backups

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
